# debai-250/network/udev/70-skillfish-wol.rules
# Re-apply Wake-on-LAN whenever a network interface is (re-)added — some
# drivers reset the WoL flag on link renegotiation, cable replug, or module
# reload, so the boot-time systemd unit alone isn't always enough.
ACTION=="add", SUBSYSTEM=="net", RUN+="/usr/local/bin/skillfish-enable-wol"
