# system/etc/systemd/system/docker.service.d/skillfish-acl.conf
[Service]
ExecStartPost=/bin/sh -c 'sleep 2; /usr/bin/setfacl -m u:skillfish:rw /var/run/docker.sock || true'
