# system/etc/systemd/system/skillfish-thermal-guard.service
[Unit]
Description=SkillFishOS thermal guard
[Service]
ExecStart=/usr/local/bin/skillfish-thermal-guard.sh
Restart=always
[Install]
WantedBy=multi-user.target
