# system/etc/systemd/system/skillfish-dp-hotswap.service
[Unit]
Description=SkillFishOS - hot-swap monitor DP (BC-250 HPD fix, EDID-watch silenzioso)
After=systemd-udev-settle.service
Before=display-manager.service

[Service]
Type=simple
ExecStart=/usr/local/bin/skillfish-dp-hotswap.sh
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
